# host/drivers/InVolFlt/linux/user_space/initrd/dracut/95involflt/99_involflt.rules
SUBSYSTEM!="block" GOTO="end_inmrule"
ACTION=="add", RUN+="/bin/sh /sbin/involflt_init.sh flt_stack %k"
LABEL="end_inmrule"
