# host/drivers/InVolFlt/linux/user_space/initrd/initramfs-tools/involflt/99_involflt.rules
SUBSYSTEM!="block" GOTO="end_inmrule"
ACTION=="add", RUN+="/bin/sh /bin/involflt_init.sh %k"
LABEL="end_inmrule"
